# desktop/.config/autostart/emacs.desktop
[Desktop Entry]
Name=Emacs Daemon
Exec=emacs --daemon
Icon=emacs
Type=Application
Terminal=false