# desktop/.config/autostart/kdeconnect.desktop
[Desktop Entry]
Name=kdeconnect daemon
Exec=/usr/lib/kdeconnectd
Type=Application
